# scripts/podman/theiaos.container.in
# TheiaOS gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into theiaos's ~/.config/containers/systemd/
# {{THEIAOS_HOME}} is replaced at install time.

[Unit]
Description=TheiaOS gateway (rootless Podman)

[Container]
Image=theiaos:local
ContainerName=theiaos
UserNS=keep-id
Volume={{THEIAOS_HOME}}/.theiaos:/home/node/.theiaos
EnvironmentFile={{THEIAOS_HOME}}/.theiaos/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
